# twal/crontab
*/1 * * * *  /root/twal/main.py >> /tmp/logs  2>&1
